# config_overlay/includes.chroot/etc/xdg/autostart/conky.desktop
[Desktop Entry]
Name=Conky
Exec=/usr/bin/conky
Terminal=false
Type=Application
