# config_overlay/includes.chroot/etc/xdg/autostart/conky.desktop
[Desktop Entry]
Name=Conky
Exec=/usr/bin/conky --pause=3 -c /usr/share/liox-config/conkyrc
Terminal=false
Type=Application
